# etc/udev/rules.d/99-grive-sync.rules
KERNEL=="eth*", RUN+="/usr/local/grive-sync/udev_program.sh" 
KERNEL=="wlan*", RUN+="/usr/local/grive-sync/udev_program.sh" 
