# etc/udev/rules.d/99-grive-sync.rules
KERNEL=="eth*", RUN+="/usr/local/grive-sync/udev_program.sh" 
KERNEL=="wlan*", RUN+="/usr/local/grive-sync/udev_program.sh" 
KERNEL=="usb_host", RUN+="/usr/local/grive-sync/udev_usb_program.sh"
KERNEL=="usb_plug", RUN+="/usr/local/grive-sync/udev_usb_program.sh"
